# deploy/systemd/trackerstream-ingest.timer
[Unit]
Description=Periodic trackerstream ingest (pick up corpus additions)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
