# deploy/systemd/trackerstream-backup.timer
[Unit]
Description=Daily trackerstream backup (catalog + pinset, prune old, DO volume snapshot)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
