# deploy/systemd/trackerstream-ingest.service
[Unit]
Description=trackerstream corpus ingest (build CID-DAGs, pin, catalog) — re-runnable
# The live master is tsnode (trackerstream-node). Must NOT depend on the kubo warm
# standby (trackerstream-ipfs): kubo Conflicts= tsnode, so requiring it here would
# evict tsnode on every ingest run — this is exactly what silently rolled prod back
# to kubo at the first nightly ingest after cutover.
After=trackerstream-node.service
Requires=trackerstream-node.service

[Service]
Type=oneshot
User=trackerstream
Group=trackerstream
EnvironmentFile=/etc/trackerstream/server.env
WorkingDirectory=/opt/trackerstream/apps/server
# Incremental: skips modules already in the catalog, so re-runs are cheap.
ExecStart=/usr/bin/env node --experimental-sqlite bin/ingest.ts
TimeoutStartSec=0
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
