# deploy/systemd/trackerstream-node.service
[Unit]
Description=trackerstream master node (tsnode) — custom go-libp2p+boxo: DHT bootstrap + Bitswap seeder + IPNS publisher + relay
After=network-online.target
Wants=network-online.target
# Mutually exclusive with the kubo warm-standby: only one master at a time (both bind
# the RPC :5001 the ingest/metrics expect; swarm ports differ — tsnode :5478, kubo :4001).
Conflicts=trackerstream-ipfs.service

[Service]
Type=simple
User=trackerstream
Group=trackerstream
# tsnode's own repo (datastore + identity.key + keystore). Distinct from the kubo
# warm-standby's /srv/trackerstream/ipfs, which is preserved for rollback.
Environment=TS_ROLE=server
Environment=TS_REPO=/srv/trackerstream/tsnode
# Non-default swarm port (was :4001) — private overlay, sheds cached public-IPFS
# scanners. Keep in sync with packages/config LIBP2P_SWARM_PORT + the client bootstrap.
Environment=TS_SWARM_PORT=5478
# Kubo-compatible RPC on the same loopback addr the ingest + API expect (KUBO_API).
Environment=TS_RPC=127.0.0.1:5001
ExecStart=/usr/local/bin/tsnode --role server
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
