# deploy/systemd/trackerstream-api.service
[Unit]
Description=trackerstream control-plane HTTP API (catalog/search)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trackerstream
Group=trackerstream
EnvironmentFile=/etc/trackerstream/server.env
WorkingDirectory=/opt/trackerstream/apps/server
ExecStart=/usr/bin/env node --experimental-sqlite bin/serve.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
